# bucket_57/dbus/files/dbus.in
#!/bin/sh
#
# PROVIDE: dbus
# REQUIRE: DAEMON ldconfig
#
# Add the following lines to /etc/rc.conf to enable the D-BUS messaging system:
#
# dbus_enable="YES"
#

. /etc/rc.subr

: ${dbus_enable=${gnome_enable-NO}} ${dbus_flags="--system"}

name=dbus
rcvar=dbus_enable

command="%%PREFIX%%/bin/dbus-daemon"
pidfile="%%PREFIX%%/var/run/dbus/pid"

stop_postcmd="dbus_poststop"

dbus_poststop()
{
    rm -f $pidfile
}

load_rc_config ${name}
run_rc_command "$1"
